# supervisord1.conf
; supervisor config file

[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)
chmod=0700                       ; sockef file mode (default 0700)

[inet_http_server] 
port=0.0.0.0:9001

[supervisord]
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
childlogdir=/var/log/supervisor            ; ('AUTO' child log dir, default $TEMP)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

; The [include] section can just contain the "files" setting.  This
; setting can list multiple files (separated by whitespace or
; newlines).  It can also contain wildcards.  The filenames are
; interpreted as relative to this file.  Included files *cannot*
; include files themselves.

[program:gost]
command=/root/app/gost/gost -C /root/app/gost/001.yaml
autostart=true
startsecs=1
autorestart=true
startretries=10
user=root
priority=999
redirect_stderr=true
stdout_logfile=/root/app/gost/out2_gost.log
stopasgroup=false
killasgroup=false

[program:GetReceive001]
directory=/root/app/GetReceive/
command=/root/app/GetReceive/001.sh
autostart=false
startsecs=1
autorestart=true
startretries=999
user=root
priority=999
redirect_stderr=true
stdout_logfile=/root/app/GetReceive/out_GetReceive.log
stopasgroup=false
killasgroup=false

[program:GetReceive002]
directory=/root/app/GetReceive2/
command=/root/app/GetReceive2/001.sh
autostart=false
startsecs=1
autorestart=true
startretries=999
user=root
priority=999
redirect_stderr=true
stdout_logfile=/root/app/GetReceive2/out_GetReceive.log
stopasgroup=false
killasgroup=false


[include]
files = /etc/supervisor/conf.d/*.conf
